# boards/common_fs/etc/update-motd.d/09-octocitrico-banner
#!/bin/bash

TERM=linux toilet -f standard -F metal -F border OctoCitrico
echo "Running on ...."